# services/corrected_paper.service
[Unit]
Description=Corrected paper executor for CORR-V4
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/trdr
Environment=BOT_ALERT_TAG=[CORR-PAPER]
ExecStartPre=/usr/bin/mkdir -p /opt/trdr/out/corrected_paper/logs /opt/trdr/out/corrected_paper/reports /opt/trdr/state/corrected_paper
ExecStart=/usr/bin/env bash /opt/trdr/scripts/run_corrected_paper.sh
Restart=always
RestartSec=5
StandardOutput=append:/opt/trdr/out/corrected_paper/logs/systemd_stdout.log
StandardError=append:/opt/trdr/out/corrected_paper/logs/systemd_stderr.log

[Install]
WantedBy=multi-user.target
